# bradley-cam.service
[Unit]
Description=Bradley.io camera snapshot (BRIO frame grab)
After=network.target

[Service]
Type=oneshot
User=bisenbek
Group=bisenbek
SupplementaryGroups=video
WorkingDirectory=/home/bisenbek/projects/bradleyio
Environment="CAM_CACHE_DIR=/var/lib/bradley-cam"
Environment="CAM_MAXW=1280"
# systemd creates/chowns /var/lib/bradley-cam for us (User:Group, 0755).
StateDirectory=bradley-cam
ExecStart=/home/bisenbek/projects/bradleyio/scripts/cam-snap.sh
SyslogIdentifier=bradley-cam
Nice=10
